# pi_code/serverroom.service
[Unit]
Description=Server Room Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/serverroom
EnvironmentFile=/opt/serverroom/.env
ExecStart=/usr/bin/python3 /opt/serverroom/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
# Ensure GPIO cleanup on stop
KillMode=process

[Install]
WantedBy=multi-user.target
